# smartmemory_app/data/ai.smartmemory.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.smartmemory.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>{PYTHON_PATH}</string>
        <string>-c</string>
        <string>from smartmemory_app.viewer_server import main; main(port={DAEMON_PORT}, open_browser=False)</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>{DATA_DIR}/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>{DATA_DIR}/daemon.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{BIN_DIR}:/usr/local/bin:/usr/bin:/bin</string>
        <key>SMARTMEMORY_DATA_DIR</key>
        <string>{DATA_DIR}</string>
    </dict>
</dict>
</plist>
